# rootfs/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/command/with-contenv bash

# Start noVNC proxy
exec s6-setuidgid $(id -nu $PUID) /usr/share/novnc/utils/novnc_proxy --listen ${NOVNC_PORT}
