# rootfs/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/command/with-contenv bash

# Print all relevant variables and command components for debugging
echo "Environment variables:"
echo "VNC_PORT=${VNC_PORT}"
echo "NOVNC_PORT=${NOVNC_PORT}"
echo "PUID=${PUID}"
echo "User resolution:"
echo "id -nu $PUID = $(id -nu $PUID)"
echo "Full command that will be executed:"
echo "s6-setuidgid $(id -nu $PUID) /usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080 --web /usr/share/novnc"

# Start noVNC proxy with explicit websockify options
exec s6-setuidgid $(id -nu $PUID) /usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080 --web /usr/share/novnc
