# rootfs/etc/s6-overlay/s6-rc.d/svc-soulseek/run
#!/command/with-contenv bash

export HOME=/data
# Default umask settings
umask ${UMASK}

# Start Soulseek
exec s6-setuidgid $(id -nu $PUID) /app/SoulseekQt
