# rootfs/etc/s6-overlay/s6-rc.d/svc-soulseek/run
#!/command/with-contenv bash

export HOME=/data

# Start Soulseek
exec s6-setuidgid soulseek /app/SoulseekQt
