# rootfs/etc/s6-overlay/s6-rc.d/svc-tigervnc/run
#!/command/with-contenv bash

# Check if the VNCPWD_FILE variable is set and file exists
if [ -n "$VNCPWD_FILE" ] && [ -f "$VNCPWD_FILE" ]; then
  VNCPWD=$(cat "$VNCPWD_FILE")
fi

# Start Tigervnc with configurable port
if [ -z "${VNCPWD}" ]; then
  exec s6-setuidgid "$(id -nu "$PUID")" Xtigervnc \
    -desktop soulseek -auth /tmp/.Xauthority -rfbport "${VNC_PORT}" \
    -nopn -quiet -AlwaysShared -SecurityTypes None :1
else
  (umask 077 && echo "${VNCPWD}" | vncpasswd -f > /tmp/passwd) || { echo "vncpasswd failed" >&2; exit 1; }
  chown "$(id -nu "$PUID")" /tmp/passwd
  exec s6-setuidgid "$(id -nu "$PUID")" Xtigervnc \
    -desktop soulseek -auth /tmp/.Xauthority -rfbport "${VNC_PORT}" \
    -nopn -rfbauth /tmp/passwd -quiet -AlwaysShared :1
fi
