# rootfs/etc/s6-overlay/s6-rc.d/svc-tigervnc/run
#!/command/with-contenv bash

# Check if the VNCPWD_FILE variable is set — fail closed if configured but unreadable/empty
if [ -n "$VNCPWD_FILE" ]; then
  [ -r "$VNCPWD_FILE" ] || { echo "VNCPWD_FILE is not readable: $VNCPWD_FILE" >&2; exit 1; }
  VNCPWD=$(cat "$VNCPWD_FILE")
  [ -n "$VNCPWD" ] || { echo "VNCPWD_FILE is empty: $VNCPWD_FILE" >&2; exit 1; }
fi

# Start Tigervnc with configurable port
if [ -z "${VNCPWD}" ]; then
  exec s6-setuidgid "$(id -nu "$PUID")" Xtigervnc \
    -desktop soulseek -auth /tmp/.Xauthority -rfbport "${VNC_PORT}" \
    -nopn -quiet -AlwaysShared -SecurityTypes None :1
else
  [ "${#VNCPWD}" -gt 8 ] && echo "WARNING: VNC auth only uses the first 8 characters of VNCPWD" >&2
  (umask 077 && printf '%s\n' "${VNCPWD}" | vncpasswd -f > /tmp/passwd) || { echo "vncpasswd failed" >&2; exit 1; }
  chown "$(id -nu "$PUID")" /tmp/passwd
  exec s6-setuidgid "$(id -nu "$PUID")" Xtigervnc \
    -desktop soulseek -auth /tmp/.Xauthority -rfbport "${VNC_PORT}" \
    -nopn -rfbauth /tmp/passwd -quiet -AlwaysShared :1
fi
